# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:node]
command=node /app/backend-node/server.js
directory=/app
autostart=true
autorestart=true
environment=PORT="%(ENV_PORT_NODE)s",DATA_DIR="%(ENV_DATA_DIR)s",PYTHON_API_URL="http://127.0.0.1:%(ENV_PORT_PYTHON)s",UNRAID_1_NAME="%(ENV_UNRAID_1_NAME)s",UNRAID_1_HOST="%(ENV_UNRAID_1_HOST)s",UNRAID_1_API_KEY="%(ENV_UNRAID_1_API_KEY)s",UNRAID_2_NAME="%(ENV_UNRAID_2_NAME)s",UNRAID_2_HOST="%(ENV_UNRAID_2_HOST)s",UNRAID_2_API_KEY="%(ENV_UNRAID_2_API_KEY)s",UNRAID_3_NAME="%(ENV_UNRAID_3_NAME)s",UNRAID_3_HOST="%(ENV_UNRAID_3_HOST)s",UNRAID_3_API_KEY="%(ENV_UNRAID_3_API_KEY)s",POLL_INTERVAL="%(ENV_POLL_INTERVAL)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:python]
command=uvicorn main:app --host 127.0.0.1 --port %(ENV_PORT_PYTHON)s --workers 2
directory=/app/backend-python
autostart=true
autorestart=true
environment=PORT="%(ENV_PORT_PYTHON)s",DATA_DIR="%(ENV_DATA_DIR)s",HOST_PROC="%(ENV_HOST_PROC)s",HOST_SYS="%(ENV_HOST_SYS)s",ADGUARD_URL="%(ENV_ADGUARD_URL)s",ADGUARD_USER="%(ENV_ADGUARD_USER)s",ADGUARD_PASSWORD="%(ENV_ADGUARD_PASSWORD)s",CROWDSEC_URL="%(ENV_CROWDSEC_URL)s",CROWDSEC_API_KEY="%(ENV_CROWDSEC_API_KEY)s",OPENWEATHER_KEY="%(ENV_OPENWEATHER_KEY)s",ICAL_URLS="%(ENV_ICAL_URLS)s",CALDAV_URL="%(ENV_CALDAV_URL)s",CALDAV_USER="%(ENV_CALDAV_USER)s",CALDAV_PASSWORD="%(ENV_CALDAV_PASSWORD)s",TZ="%(ENV_TZ)s"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
